# nightshift/deploy/nightshift-scrape.service
# /etc/systemd/system/nightshift-scrape.service
# Refreshes the planes-for-sale inventory daily (ingest + alert match-and-send).
[Unit]
Description=ClubHanger Night Shift daily inventory scrape
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
TimeoutStartSec=2700
ExecStart=/usr/bin/docker run --rm \
  --name clubhanger-nightshift-scrape \
  --memory=2g --cpus=1 \
  -e NS_STATE_DIR=/home/night/state \
  -e TZ=America/Los_Angeles \
  -v /opt/nightshift/repo:/app \
  -v /opt/nightshift/claude:/home/night/.claude \
  -v /opt/nightshift/state:/home/night/state \
  -v /opt/nightshift/ssh:/home/night/.ssh:ro \
  --entrypoint /app/nightshift/bin/run-scrape.sh \
  clubhanger-nightshift:latest

[Install]
WantedBy=multi-user.target
